# src/etc/supervisord.d/ratelimitter.ini
[program:ratelimitter]
command=/usr/bin/python3  ratelimitter_server.py
directory=/opt/ratelimitter
user=rluser
autostart=true
autorestart=true

